# config/armshell.conf
[program:armshell]
command=/home/Karlson/armshell/venv/bin/gunicorn config.wsgi.application -c /home/Karlson/armshell/config/gunicorn.conf.py
directory=/home/Karlson/armshell
user=Karlson
autorestart=true
redirect_stderr=true
stdout_logfile=/home/Karlson/armshell/logs/debug.log